# gswitch.service
[Unit]
Description=Automatic GPU switching service

[Service]
ExecStart=/usr/bin/gswitch boot

[Install]
WantedBy=multi-user.target
